# pipeline/run/super_structure_parity.service
[Unit]
Description=Super Structure Parity Check (3-way ledger vs signals vs UI)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/kemal/futures
ExecStart=/usr/bin/python3 -u /home/kemal/futures/pipeline/live/parity_super_structure.py --telegram auto
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/kemal/futures/data/Live/super_structure_parity.log
StandardError=append:/home/kemal/futures/data/Live/super_structure_parity.log
